# etc/local.d/10_disable_intel_pstate_turbo_mode.start
start/10_disable_intel_pstate_turbo_mode.start