# etc/local.d/11_disable_cpu_hyperthreading.start
#!/bin/sh
script_directory_path="${0%/*}"
script_name="${0##*/}"

/usr/bin/logger --tag "disable cpu hyperthreading" --id --stderr "${script_directory_path}/${script_name}: executed"

# disable cpu hyperthreading on cpu 1, 3, 5 and 7
for cpu_number in {1..7..2}
do
    echo "0" > "/sys/devices/system/cpu/cpu${cpu_number}/online"
done
