# etc/local.d/11_disable_cpu_hyperthreading.start
start/11_disable_cpu_hyperthreading.start