# etc/local.d/12_scaling_max_freq.start
#!/bin/sh
script_directory_path="${0%/*}"
script_name="${0##*/}"

/usr/bin/logger --tag "adjust scaling_max_freq" --id --stderr "${script_directory_path}/${script_name}: executed"

for policy in "/sys/devices/system/cpu/cpufreq/policy"*
do
    # downclock all cpu cores to 2.4ghz (default)
    /bin/echo "2400000" > "${policy}/scaling_max_freq"
done
